# files/supervisord.conf
[supervisord]
nodaemon=true

[program:sinopia]
command=/var/opt/sinopia/bin/sinopia

[program:gwp]
command=/var/opt/gitlab-webhook-publish/bin/gwp

[program:gitlab]
command=/usr/local/bin/wrapper